# SFSCextractor@.service
[Unit]
Description=Start SFSCextractor

[Service]
Type=simple
ExecStart=/%h/bin/SFSCextractor.sh

[Install]
WantedBy=default.target
