# SFSCextractor@.service
[Unit]
Description=Start SFSCextractor instance

[Service]
Type=simple
ExecStart=/usr/bin/SFSCextractor.sh

[Install]
WantedBy=default.target
